# scripts/completion-service/ag-completion-service.socket
[Unit]
Description=ag CLI completion service socket

[Socket]
ListenStream=%t/ag-completion.sock

[Install]
WantedBy=sockets.target
